# systemd/piclock-renderer.service
[Unit]
Description=Pi Klydo Clock renderer
Documentation=file:/opt/piclock/README.md
# Start as early as the filesystem allows. We deliberately do NOT reference
# network-online.target, so the first frame is never gated on Wi-Fi.
After=local-fs.target systemd-user-sessions.service
Wants=local-fs.target
# Take over the console; keep getty off tty1.
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
# install.sh rewrites User/Group to PICLOCK_USER, defaulting to the sudo user.
User=piclock
Group=piclock
# DRM/KMS + input access.
SupplementaryGroups=video render input tty
StateDirectory=piclock
WorkingDirectory=/opt/piclock/src
Environment=SDL_VIDEODRIVER=kmsdrm
Environment=PYTHONPATH=/opt/piclock/src
# Display, layout, design folders, and rotation are loaded from clock.json.
ExecStart=/opt/piclock/.venv/bin/python -m piclock.cli
Restart=always
RestartSec=2
# Give the process the VT so SDL can become DRM master.
StandardInput=tty
StandardOutput=journal
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
